# konar@.timer
[Unit]
Description=scheduled runescape hiscore download

[Install]
WantedBy=timers.target

[Timer]
OnUnitInactiveSec=1h
RandomizedDelaySec=5m
