# konar@.timer
[Unit]
Description=scheduled runescape hiscore download

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5m
